# backend/dotfiles/scanner.service
# =============================================================================
# CVE Scan API — unit systemd
# Skopiuj do /etc/systemd/system/cve-scan-api.service
# Po skopiowaniu: sudo systemctl daemon-reload && sudo systemctl enable --now cve-scan-api
# =============================================================================

[Unit]
Description=CVE Scan API — debsecan + NVD vulnerability scanner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

User=cveapi
Group=cveapi

WorkingDirectory=/opt/cve-scan-api

# Plik ze zmiennymi środowiskowymi (chmod 600, właściciel root)
EnvironmentFile=/etc/cve-scan-api/scanner.env

# Uruchomienie przez virtualenv (zalecane) lub systemowy Python3
ExecStart=/opt/cve-scan-api/.venv/bin/python /opt/cve-scan-api/main.py
# ExecStart=/usr/bin/python3 /opt/cve-scan-api/main.py

Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30

# --- Hardening ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

# Katalog cache NVD musi być zapisywalny (zgodny z CVE_CACHE_PATH w scanner.env)
ReadWritePaths=/var/cache/cve-scan-api

# --- Logi przez journald ---
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
